# fakessh/systemd/fakessh.service
[Unit]
Description=Fake SSH server
After=network.target

[Service]
ExecStart=/usr/bin/fakessh

[Install]
WantedBy=multi-user.target
